# settings/pigarage_project.conf
[program:pigarage]
command = python3 pigarage.py
directory = /pigarage
autostart = true
autorestart = true
environment=PIGARAGE_SETTINGS="/settings/settings.cfg"